# ldapserver/rc.local
#! /bin/bash

cd /etc/pam.d
sed -i 's/^.*sss.*$//g' *

dhclient -v {{ifconfig_output.stdout}}
